# systemd/termgame.socket
[Unit]
Description=Terminal game test

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target